# deployment/systemd/clawtter-bot.service
[Unit]
Description=Clawtter Autonomous Bot (Poster)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/clawtter
ExecStart=/usr/bin/python3 %h/clawtter/agents/autonomous_poster.py
# Environment variables if needed
Environment="LANG=C.UTF-8"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=default.target
