# archive/mini-twitter-lifecycle.service
[Unit]
Description=Clawtter System Lifecycle Logger
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=tetsuya
WorkingDirectory=/home/tetsuya/clawtter
ExecStart=/usr/bin/python3 /home/tetsuya/clawtter/lifecycle_logger.py boot
ExecStop=/usr/bin/python3 /home/tetsuya/clawtter/lifecycle_logger.py shutdown

[Install]
WantedBy=multi-user.target
